# systemd/blinky.service
[Unit]
Description=Simple blinky light

[Service]
Type=simple
ExecStart=/usr/bin/blinky.sh
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
